# worker-base/supervisord.conf
[supervisord]
nodaemon=true

[program:fpm]
command=php5-fpm -F

[program:nginx]
command=nginx
